# tf_pack/luci/applications/luci-app-tfortis-ups/root/etc/uci-defaults/80_tfortis-ups
#!/bin/sh

touch /etc/config/tfortis_ups


return 0
